# web_application_old/health-monitor.service
[Unit]
Description=Health Monitor for Pesticide Search Application
After=network.target pesticide-search.service
Requires=pesticide-search.service

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/pesticide-search
Environment=PATH=/home/ubuntu/pesticide-search/venv/bin
ExecStart=/home/ubuntu/pesticide-search/venv/bin/python health_monitor.py
Restart=always
RestartSec=30
StandardOutput=journal
StandardError=journal
SyslogIdentifier=health-monitor

[Install]
WantedBy=multi-user.target




